# dockerfile/apollo/supervisord.conf
[supervisord]  
nodaemon=true  
  
[program:config]  
command=/runconfig.sh
  
[program:admin]  
command=/runadmin.sh

[program:portal]  
command=/runportal.sh
